# service/autoscaler.service
[Unit]
Description=KVM autoscaler

[Service]
ExecStart=autoscaler

[Install]
WantedBy=multi-user.target
